# configs/supervisor.d/strongswan.ini
[program:strongswan]
command=/usr/bin/ipsec start --nofork
redirect_stderr=true
stdout_logfile=/var/log/supervisor/strongswan.log
