# configs/supervisor.d/swanctl.ini
[program:swanctl]
command=/usr/lib/strongswan/charon
redirect_stderr=true
stdout_logfile=/var/log/supervisor/swanctl.log
